# workspace/docs/edict-dashboard.service
[Unit]
Description=Edict Dashboard Server
After=network.target openclaw-gateway.service
Wants=openclaw-gateway.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/projects/edict
ExecStart=/usr/bin/python3 /root/.openclaw/workspace/projects/edict/dashboard/server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
